# systemd/emailguard-backend.service
[Unit]
Description=EmailGuard Django Backend (Gunicorn)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data

WorkingDirectory=/opt/emailguard/backend

# Load environment
EnvironmentFile=/opt/emailguard/backend/.env
Environment="DJANGO_SETTINGS_MODULE=config.settings.local"
Environment="PATH=/opt/emailguard/backend/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Gunicorn — 3 workers, binds to port 9002
ExecStart=/opt/emailguard/backend/venv/bin/gunicorn \
    config.wsgi:application \
    --bind 0.0.0.0:9002 \
    --workers 3 \
    --worker-class sync \
    --timeout 120 \
    --keep-alive 5 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --log-level info \
    --access-logfile /var/log/emailguard/backend-access.log \
    --error-logfile /var/log/emailguard/backend-error.log

ExecReload=/bin/kill -s HUP $MAINPID

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=emailguard-backend

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
